# dist/pkg-deb/amd64/DEBIAN/prerm
systemctl stop sensible
systemctl disable sensible
systemctl daemon-reload
rm -rf /var/log/sensible/*